# deploy/bank-discrepancies.service
[Unit]
Description=Bank Discrepancies (הפקת קבלות) — Flask/gunicorn
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Bank-discrepancies

# Real secrets live OUTSIDE the repo, in a central env folder (chmod 600).
# The app reads this path first (see backend/server.py env resolution).
Environment=BANK_ENV_FILE=/home/ubuntu/env/.env
Environment=PATH=/home/ubuntu/Bank-discrepancies/.venv/bin:/usr/bin

# gunicorn serves the Flask app object `app` from backend/server.py.
# Bind to localhost only — nginx is the public-facing reverse proxy.
ExecStart=/home/ubuntu/Bank-discrepancies/.venv/bin/gunicorn \
    --chdir /home/ubuntu/Bank-discrepancies/backend \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --threads 4 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    server:app

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
